# systemd/user/battery-notify.service
[Unit]
Description=Battery notification check

[Service]
ExecStart=%h/.local/bin/battery-notify.sh
